# buffer-sharer-middleware/buffer-sharer-middleware.service
[Unit]
Description=Buffer Sharer Middleware Server
After=network.target
Documentation=https://github.com/buffer-sharer/middleware

[Service]
Type=simple
User=nobody
Group=nogroup
ExecStart=/usr/local/bin/buffer-sharer-middleware -port 8080
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes

[Install]
WantedBy=multi-user.target
